# supervisord.conf
[supervisord]
nodaemon=true

[program:heartbeat]
command=python3 heartbeat_mailing.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/heartbeat_stdout.log
stderr_logfile=/var/log/heartbeat_stderr.log
redirect_stderr=true

[program:consumer]
command=python3 consumer_mailing.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/consumer_stdout.log
stderr_logfile=/var/log/consumer_stderr.log
redirect_stderr=true